# scripts/profile.d/global-constitution.uss.sh
echo test with alias
alias aaa="echo moo"
